# Sandboxer/sandboxer.service
[Unit]
Description=Sandboxer software isolation suite

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/true
ExecStop=/usr/bin/sandboxer-stop-all 10 10

[Install]
WantedBy=multi-user.target
